# upstart.conf
description 'es sync'
start on runlevel [2345]
stop on runlevel [06]

respawn
normal exit 0

chdir <PATH_TO_CONFIG>

script
    es-sync config.yaml
end script